# cron.d/lynis-warnings-email
## Lynis: run system audit and email warnings
20 4 * * * root /root/scripts/lynis-warnings-email.sh
