# cron.d/lynis-warnings-email
# Lynis: Conduct a system audit and dispatch email notifications for any detected warnings.
20 4 * * * root /root/scripts/lynis-warnings-email.sh
