# cron.d/lynis-warnings-email
## Lynis: Perform a system audit and send email notifications for any warnings detected.
20 4 * * * root /root/scripts/lynis-warnings-email.sh
