# cron.d/needrestart-services-flag
## Create the reboot-required file if some services need to be restarted
0 */4 * * * root /root/scripts/needrestart-services-flag.sh
