# cron.d/needrestart-services-flag
## Create the "reboot-required" file if there are services that require a restart.
0 */4 * * * root /root/scripts/needrestart-services-flag.sh
