# cron.d/needrestart-services-flag
# Generate the 'reboot-required' file if there are services that necessitate a restart.
0 */4 * * * root /root/scripts/needrestart-services-flag.sh
