# cron.d/reboot-required-email
## Check if reboot-required flag is set and email notification
00 05 * * * root /root/scripts/reboot-required-email.sh
