# cron.d/reboot-required-email
## Verify if the 'reboot-required' flag is set and dispatch an email notification.
00 05 * * * root /root/scripts/reboot-required-email.sh
