# cron.d/reboot-required-email
## Check if the "reboot-required" flag is set and send an email notification.
00 05 * * * root /root/scripts/reboot-required-email.sh
